# deploy/systemd/akfa-backend-watchdog.timer
[Unit]
Description=Run AKFA backend watchdog every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=60s
AccuracySec=10s
Persistent=true

[Install]
WantedBy=timers.target
